# app/third_party/jupyterhub/conf/systemd/jupyter_notebook.service
[Unit]
Description=Job that runs the jupyter_notebook daemon

[Service]
User=${JUPYTERHUB_SERVICE_USER}
Group=${JUPYTERHUB_SERVICE_GROUP}
Environment=VIRTUAL_ENV=${JUPYTERHUB_VENV}
Environment=PYTHONPATH=${JUPYTERHUB_VENV}
WorkingDirectory=${JUPYTERHUB_VENV}
ExecStart=${JUPYTERHUB_VENV}/bin/jupyter notebook \
  --NotebookApp.notebook_dir='${JUPYTERHUB_NOTEBOOK_DIR}' \
  --NotebookApp.ip='${JUPYTERHUB_IP}' \
  --NotebookApp.port='${JUPYTERHUB_PORT}' \
  --Session.username='${JUPYTERHUB_USERNAME}' \
  --NotebookApp.password='${JUPYTERHUB_PASSWORD}' \
  --NotebookApp.password_required=True \
  --NotebookApp.allow_remote_access=True \
  --NotebookApp.iopub_data_rate_limit='2147483647' \
  --no-browser \
  --NotebookApp.open_browser=False
#  --NotebookApp.certfile='/etc/letsencrypt/live/REPLACE_WITH_DNS_NAME/fullchain.pem' \
#   --NotebookApp.keyfile='/etc/letsencrypt/live/REPLACE_WITH_DNS_NAME/privkey.pem' \

[Install]
WantedBy=multi-user.target
